# supervisord.conf
[supervisord]
nodaemon=true
[program:plex]
command=/bin/sh -c '/usr/lib/plexmediaserver/Plex\ Media\ Server'
