# supervisord.conf
[supervisord]
nodaemon=true
[program:plex]
command=/usr/sbin/start_pms
